# config/supervisor/mod.conf
[program:mod]
command=/usr/local/bin/odr-dabmod /home/odr/config/mod.ini
autostart=false
autorestart=true
priority=1
stderr_logfile=/var/log/supervisor/mod.log
stdout_logfile=/var/log/supervisor/mod.log
